# libratbag/ghostcatd/ghostcatd.service.in
[Unit]
Description=Daemon to introspect and modify configurable mice

[Service]
Type=simple
BusName=org.freedesktop.ghostcat1
ExecStart=@sbindir@/ghostcatd
Restart=on-failure

[Install]
Alias=dbus-org.freedesktop.ghostcat1.service
